# config/systemd/notmuch.service
[Unit]
Description=Mail indexer cleanup service

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch compress

[Install]
WantedBy=default.target
